# src/container/kube-manager/contrail-kube-manager.ini
[program:tungsten-kube-manager]
command=/usr/bin/tungsten-kube-manager -c /etc/tungsten/tungsten-kubernetes.conf
priority=460
autostart=true
killasgroup=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/tungsten/tungsten-kube-manager-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=tungsten
