# src/config/api-server/contrail-api.ini
[program:tungsten-api]
command=/usr/bin/tungsten-api --conf_file /etc/tungsten/tungsten-api.conf --conf_file /etc/tungsten/tungsten-keystone-auth.conf --worker_id %(process_num)s
numprocs=1
process_name=%(process_num)s
redirect_stderr=true
stdout_logfile= /var/log/tungsten/tungsten-api-%(process_num)s-stdout.log
stderr_logfile=/dev/null
priority=440
autostart=true
killasgroup=true
stopsignal=TERM
exitcodes=0 ; 'expected' exit codes for process (default 0,2)
user=tungsten
