# src/dns/contrail-named.ini
[program:tungsten-named]
command=/usr/bin/authbind /usr/bin/tungsten-named -f -c /etc/tungsten/dns/tungsten-named.conf
user=tungsten
priority=520
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/tungsten/tungsten-named-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=tungsten
